# spell-check.ini
[program:spell-check]
directory=/var/spell-check
command=/var/spell-check/main.py
autostart=true
autorestart=true
startretries=3
