# spell-check.ini
[program:spell-check]
directory=/var/spell-check
command=/var/spell-check/venv/bin/python main.py
autostart=true
autorestart=true
stdout_logfile=/var/log/spell-check.log
startretries=3
